# hypr-py-light.service
[Unit]
Description=Hypr-Py-Light Screen Temperature Manager
Documentation=https://github.com/IamAndelib/Auto-blue-light-filter
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/hypr-py-light
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
